# hostplumber/pkg/ovs-docker/ovs.conf
[program:ovs]
command=/usr/local/bin/start-ovs.sh
